# tikrtrakrd.service
[Unit]
Description=daemon for tikrtrakr

[Service]
Type=forking
User=root
ExecStart=tikrtrakrd XRPUSDT 60

[Install]
WantedBy=multi-user.target
